# ansible-idempotent-automation/templates/app.service.j2
[Unit]
Description=My Application v{{ app_version }}
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/app-{{ app_version }}
ExecStart=/usr/local/bin/myapp-{{ app_version }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
